# providers/azure/native/rc.local
export PATH=$PATH:/usr/local/sbin:/usr/local/bin

sync-metadata
