# providers/rc.local
export PATH=$PATH:/usr/local/sbin:/usr/local/bin

echo "SSH server fingerprints:"
find /etc/ssh/ -name 'ssh_host_*.pub' -exec ssh-keygen -l -f '{}' ';'
find /etc/ssh/ -name 'ssh_host_*.pub' -exec ssh-keygen -l -E md5 -f '{}' ';'
echo

if [ -f /etc/rc.local_custom ]; then
	echo "running rc.local_custom"
	sh /etc/rc.local_custom || exit $?
fi

